# gmail_newsletter_tts/com.easyuji.newsletter-tts.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.easyuji.newsletter-tts</string>

    <key>ProgramArguments</key>
    <array>
        <string>__REPO_DIR__/gmail_newsletter_tts/run.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__REPO_DIR__/gmail_newsletter_tts</string>

    <!-- 毎朝7時に実行 -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>7</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <!-- スリープ中だった場合、起動後に実行 -->
    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>__HOME__/newsletter-tts.log</string>

    <key>StandardErrorPath</key>
    <string>__HOME__/newsletter-tts-error.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>__HOME__</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>LANG</key>
        <string>ja_JP.UTF-8</string>
    </dict>
</dict>
</plist>
